# crontab
* * * * * bundle exec ruby sample_data.rb
